# configuration/etc/init/ritual.conf
start on filesystem
exec /home/samdm/Projects/pewpewthespells-com/ritual /home/samdm/Projects/WebsiteBackup/sitemap.toml
